# scripts/backup-config/hjtpx-verify-backups.service
[Unit]
Description=HJTPX Backup Verification Service
After=network.target postgresql.service

[Service]
Type=oneshot
User=postgres
WorkingDirectory=/workspace/hjtpx
ExecStart=/workspace/hjtpx/scripts/verify-backup.sh latest
StandardOutput=append:/workspace/hjtpx/logs/backup_verify.log
StandardError=append:/workspace/hjtpx/logs/backup_verify.log
Environment="BACKUP_DIR=/var/backups/hjtpx"

