# scripts/backup-config/hjtpx-backup-full.timer
[Unit]
Description=HJTPX Full Database Backup Timer

[Timer]
OnCalendar=Sun *-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target

